# overlays/gnome/etc/skel/.config/autostart/trash-icon.desktop
[Desktop Entry]
Type=Application
Name=trash-icon
Exec=sh -c "$HOME/bins/trash.sh -d"
