# overlays/mate-fta/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "sleep 10; conky;"
Name=Conky
Comment=Autostart conky at login
