# overlays/gnome/etc/profile.d/disable-vblank.sh
#!/bin/bash

export CLUTTER_VBLANK=none
