# overlays/sway/etc/profile.d/qt-wayland.sh
#!/bin/sh

export QT_QPA_PLATFORM="wayland"
export QT_QPA_PLATFORMTHEME=qt5ct
export QT_WAYLAND_DISABLE_WINDOWDECORATION="1"
