# overlays/nemomobile/etc/systemd/system/user@.service.d/local.conf
[Service]
Environment="QT_QPA_PLATFORM=wayland"
Environment="QT_QUICK_CONTROLS_MOBILE=1"
Environment="QT_WAYLAND_DISABLE_WINDOWDECORATION=1"
Environment="MALIIT_FORCE_DBUS_CONNECTION=1"