# overlays/gnome/usr/lib/systemd/system/remove-classic.service
[Unit]
Description=Disable Gnome Classic

[Service]
ExecStart=/usr/local/bin/remove-classic.sh

[Install]
WantedBy=multi-user.target
