# overlays/gnome/usr/lib/systemd/system/remove-classic.service
[Unit]
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/remove-classic.sh

[Install]
WantedBy=default.target
